# data/mausberry-switch.service
[Unit]
Description=Mausberry switch GPIO monitoring daemon
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/mausberry-switch
Restart=on-failure
RestartSec=5s
StartLimitBurst=3
StartLimitIntervalSec=30s

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

# GPIO access requires /dev/gpiochip*
DeviceAllow=/dev/gpiochip0 rw

[Install]
WantedBy=multi-user.target
